# s6/debian-root/etc/services.d/fancy-service/run
#!/usr/bin/with-contenv bash

python3 service.py
